# code/master/install/assets/tethys-api.service
[Unit]
Description=Tethys API
After=network.target

[Service]
ExecStart={TETHYS-PATH}/env_tethys/bin/gunicorn -c {TETHYS-PATH}/install/assets-localized/gunicorn_config_api.py tethys_api.wsgi
WorkingDirectory={TETHYS-PATH}/api/
Environment="PATH={TETHYS-PATH}/env_tethys/bin"
Environment="TETHYS_DEBUG={DEBUG-FLAG}"
Environment="TETHYS_ALLOWED_HOSTS={ALLOWED-HOSTS}"
StandardOutput=inherit
StandardError=inherit
Restart=always
#User=root

[Install]
WantedBy=multi-user.target
